# webServer/supervisor.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /home/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx
